# contrib/debian/DEBIAN/prerm
#!/bin/sh

# shutdown aa_monitor service
service aa_monitor stop

# aa_monitor shouldn't start at system startup
update-rc.d -f aa_monitor disable
update-rc.d -f aa_monitor remove

# EOF